# SenseHatMonitor.service.template
[Unit]
Description=SenseHatMonitor monitors sensors on the SenseHAT

[Service]
Type=simple
ExecStart=[PWD]

[Install]
WantedBy=multi-user.target
